# blackjack.service
[Unit]
Description=BlackJack Game
[Service]
ExecStart=/bin/bash /usr/sbin/blackjack.sh
[Install]
WantedBy=multi-user.target
